# 90-hdmi.rules
ACTION=="REMOVE", SUBSYSTEM=="misc", KERNEL=="gpiohdmi", RUN+="echo on > /sys/d$

