# 90-hdmi.rules
ACTION=="remove", SUBSYSTEM=="misc" ,KERNEL=="gpiohdmi", RUN+="/usr/bin/sudo /home/pi/ir_send.sh"

